# zerotier-gateway/rootfs/etc/services.d/gateway/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Gateway Manager stopped"