# zerotier-gateway/rootfs/etc/services.d/gateway-manager/finish
#!/usr/bin/with-contenv bashio

# Clean up iptables rules on exit
if [ -f /tmp/zerotier-gateway/network.json ]; then
    ZT_DEVICE=$(jq -r '.device' /tmp/zerotier-gateway/network.json 2>/dev/null || echo "")
    ZT_IP=$(jq -r '.ip' /tmp/zerotier-gateway/network.json 2>/dev/null || echo "")
    LOCAL_SUBNET=$(bashio::config 'local_subnet')
    DEFAULT_IF=$(ip route | grep default | awk '{print $5}' | head -n1)
    
    if [[ -n "$ZT_DEVICE" && -n "$ZT_IP" && -n "$DEFAULT_IF" ]]; then
        bashio::log.info "Cleaning up iptables rules..."
        
        # Remove our specific rules
        iptables -t nat -D POSTROUTING -s ${ZT_IP}/24 -d ${LOCAL_SUBNET} -j MASQUERADE 2>/dev/null || true
        iptables -t nat -D POSTROUTING -s ${ZT_IP}/24 -o ${DEFAULT_IF} -j MASQUERADE 2>/dev/null || true
        
        iptables -D FORWARD -i ${ZT_DEVICE} -o ${DEFAULT_IF} -j ACCEPT 2>/dev/null || true
        iptables -D FORWARD -i ${DEFAULT_IF} -o ${ZT_DEVICE} -m state --state RELATED,ESTABLISHED -j ACCEPT 2>/dev/null || true
        iptables -D FORWARD -i ${ZT_DEVICE} -d ${LOCAL_SUBNET} -j ACCEPT 2>/dev/null || true
        iptables -D FORWARD -s ${LOCAL_SUBNET} -o ${ZT_DEVICE} -j ACCEPT 2>/dev/null || true
    fi
fi