# zerotier-gateway/rootfs/etc/services.d/gateway-manager/run
#!/usr/bin/with-contenv bashio

# Run the gateway manager
exec /usr/bin/gateway-manager.sh