# zerotier-gateway/rootfs/etc/services.d/zerotier/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Zerotier service stopped"