# zerotier-gateway/rootfs/etc/services.d/zerotier/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Zerotier One daemon..."

# Start zerotier-one in foreground mode
exec zerotier-one -U